# deployment/titus-streamlit.service
[Unit]
Description=Titus Simulator Web UI (Streamlit)
After=network.target titus-simulator.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/titusSimulator-v0.1
Environment="PATH=/home/ubuntu/titusSimulator-v0.1/venv/bin:/usr/bin"
Environment="PYTHONPATH=/home/ubuntu/titusSimulator-v0.1/src"
ExecStart=/home/ubuntu/titusSimulator-v0.1/venv/bin/python -m streamlit run streamlit_ui.py --server.port 8088 --server.address 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
